# bash/gh_user_info/gh_user_info.timer
[Unit]
Description=Run gh_user_info every hour
[Timer]
OnCalender=hourly
